# HTBOX-via-forlinx/rootfs-custom-add/etc/network/if-pre-up.d/wlan0
#!/bin/sh

if [[ ${IFACE} = wlan* ]]; then
	/usr/local/sbin/wpa_supplicant -B -Dwext -i${IFACE} -c/etc/wpa_supplicant.conf
fi